# scripts/deb/prerm
#!/bin/sh
echo "**** Stopping the Adafruit WebIDE ****"
service adafruit-webide.sh stop
echo "**** Removing update-rc.d service ****"
update-rc.d -f adafruit-webide.sh remove

exit 0